# cli/cmd/earbug-gchat/earbug-gchat.service
[Unit]
Description=report listening records

[Service]
EnvironmentFile=/etc/mono/earbug-gchat/env
ExecStart=/usr/local/bin/earbug-gchat
